# debian/sonm-dwh.service
[Unit]
Description=SONM DWH


[Service]
EnvironmentFile=/etc/default/sonm-dwh
Type=simple
ExecStart=/usr/bin/sonmdwh --config=${CONFIG_PATH}

Restart=on-failure
StartLimitIntervalSec=60
StartLimitBurst=3


[Install]
WantedBy=multi-user.target
